# linux/.config/systemd/user/backups.timer
[Unit]
Description=Run backups daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
