# linux/.config/systemd/user/backups.target
[Service]
Description=Backup Target
StopWhenUnneeded=yes
