# linux/.config/systemd/user/backups.timer
[Unit]
Description=Run backups daily

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=default.target
