# linux/.config/systemd/user/backups.service
[Unit]
Description=Backups my homedir
Wants=backups.timer

[Service]
ExecStart=/home/rando/bin/backup.sh

[Install]
WantedBy=backups.target
